# scripts/auto-garment-web.service
[Unit]
Description=Auto Garment Web Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/path/to/auto-garment-web
Environment=PATH=/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/path/to/auto-garment-web/.env
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
